# files/system/common/etc/profile.d/micro.sh
# Set micro's colorscheme to 'simple', which works with both light/dark themes.

alias micro='micro --colorscheme=simple'
